# flask-helpdesk.service
[Unit]
Description=Flask Helpdesk Application
After=network.target

[Service]
Type=exec
User=yvarslavan
Group=yvarslavan
WorkingDirectory=/opt/www/its.teztour.com
Environment=PATH=/opt/www/its.teztour.com/venv/bin
Environment=FLASK_ENV=production
Environment=PYTHONPATH=/opt/www/its.teztour.com
Environment=SESSION_FILE_DIR=/tmp/flask_sessions
ExecStart=/opt/www/its.teztour.com/venv/bin/gunicorn \
          --workers 1 --threads 4 --timeout 120 \
          --bind unix:/run/gunicorn/gunicorn.sock \
          --access-logfile logs/access.log \
          --error-logfile logs/error.log \
          --log-level info \
          wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Директория для сокета создаём на старте, если ещё нет
RuntimeDirectory=gunicorn
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
